# supervisor/conf.d/gunicorn_app.conf
[program:gunicorn_app]
command = python3 /root/.pyenv/versions/3.6.3/bin/gunicorn app:app --config /app/gunicorn_conf.py
directory = /app/
user = root
